# units/plex.service.d/00-plexpass.conf
# Drop-in to enable the use of the Plex Pass version of Plex and
# the features that it adds.
[Service]
Environment='PLEX_IMAGE=timhaak/plexpass'
