# units/plex.service.d/05-data.conf
[Service]
Environment='PLEX_DATA=/data'
