# units/plex.service.d/30-net-host.conf
[Service]
Environment='DOCKER_NETWORKING=--net=host'
