# units/plex.service.d/50-portmap.conf
[Service]
Environment='DOCKER_PORTMAPS=-p 32400:32400'
